# infra/systemd/backend.service.example
[Unit]
Description=School Payments Backend
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/app/backend
EnvironmentFile=/etc/app-backend.env
ExecStart=/usr/bin/npm run start:prod
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target


